# deploy/safearound-area-tracker.service
[Unit]
Description=SafeAround Small Area Crime Tracker
After=network.target postgresql.service

[Service]
Type=simple
User=safearound
WorkingDirectory=/Users/harshsaini/Desktop/Safe-Around/crime-pipeline
ExecStart=/usr/bin/python3 area_specific_tracker.py
Restart=always
RestartSec=10
Environment="NEWS_API_KEY=your_api_key"
Environment="DB_HOST=localhost"
Environment="DB_PORT=5432"
Environment="DB_NAME=safearound_prod"
Environment="DB_USER=safearound_user"
Environment="DB_PASSWORD=password"

[Install]
WantedBy=multi-user.target
