# deploy/safearound-worker.service
[Unit]
Description=SafeAround Background Worker
Documentation=https://github.com/obsidiannnn/Safe-Around
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=safearound
Group=safearound
WorkingDirectory=/opt/safearound

ExecStart=/opt/safearound/worker

# Restart policy
Restart=always
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=3

# ---- Environment ----
EnvironmentFile=/etc/safearound/worker.env

# Fallback non-secret envs
Environment=DB_SSLMODE=disable
Environment=TZ=Asia/Kolkata

# ---- Logging ----
StandardOutput=journal
StandardError=journal
SyslogIdentifier=safearound-worker

# ---- Security hardening ----
NoNewPrivileges=yes
ProtectSystem=strict
PrivateTmp=yes
ReadWritePaths=/opt/safearound

[Install]
WantedBy=multi-user.target
